# campingroute-counter.service
[Unit]
Description=CampingRoute.app Counter Service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=root
WorkingDirectory=/home/kopi/route-planner-pro
ExecStart=/usr/bin/node /home/kopi/route-planner-pro/server/index.cjs
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=NODE_ENV=production
Environment=PORT=3001
Environment=HOST=0.0.0.0
Environment=GEOAPIFY_API_KEY=replace_me
Environment=PLACE_DATABASE_PATH=/home/kopi/route-planner-pro/places.sqlite
Environment=PLACE_INDEX_PATH=/home/kopi/route-planner-pro/place-index.json
# Optional but recommended for production builds with Geoapify map tiles:
# Environment=VITE_GEOAPIFY_MAPS_API_KEY=replace_me

# Prevent systemd from considering the service as failed if it exits with success
SuccessExitStatus=0

# Increase timeout
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
